# finnazureflaskapp.conf
[program:finn-flask-app]
directory=/home/soumya/python/finnazureflaskapp
command=python main.py
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
